# extra/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

; Next.js standalone server (frontend)
[program:web]
command=node /app/web/server.js
directory=/app/web
environment=NODE_ENV="production",PORT="8000",HOSTNAME="0.0.0.0",NEXT_TELEMETRY_DISABLED="1"
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; FastAPI backend
[program:api]
command=/app/api/.venv/bin/python /app/api/app.py
directory=/app/api
environment=PYTHONUNBUFFERED="1",PYTHONIOENCODING="utf-8"
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
